# scripts/startup.service
[Unit]
Description=Startup

[Service]
Type=simple
ExecStart=/bin/Startup.sh

[Install]
WantedBy=multi-user.target
 
